# deploy/systemd/b2b-mssql-backup.timer
[Unit]
Description=Daily B2B MSSQL backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target

